# astibot.service
[Unit]
Description=Astibot Trading Bot
Documentation=https://github.com/yourusername/astibot
After=network.target

[Service]
Type=simple
User=pi
WorkingDirectory=/home/pi/astibot
Environment=NODE_ENV=production
Environment=PORT=3000
ExecStart=/usr/bin/node server.js
Restart=on-failure
RestartSec=10
StandardOutput=syslog
StandardError=syslog
SyslogIdentifier=astibot

# Security settings
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
